# etc/systemd/system/nvidia_fix.service
[Unit]
Description=Nvidia dGPU fix

[Service]
ExecStart=/usr/bin/nvidia_fix.sh

[Install]
WantedBy=multi-user.target
